# rc.local
#!/bin/sh -e
#
# rc.local
#
/boot/retrogame &
